# deploy/intraday-watch.service
[Unit]
Description=vaxstock v2 Intraday Watch (盘中盯盘推送)
After=network-online.target stock-api.service
Wants=network-online.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/stock-reportv2/vax-stock-public
EnvironmentFile=/etc/vaxstock/vaxstock.env
Environment="PATH=/opt/stock-reportv2/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/stock-reportv2/venv/bin/python -m vaxstock.services.intraday
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
